# etc/profile.d/71less.sh
#
export LESS='-R -M -i --shift=1'
export LESSCOLOR='always'
